# deploy/hermes-dashboard-monitor.service
[Unit]
Description=Hermes Dashboard State Monitor (agent.log real-time)
After=network.target hermes-dashboard.service
Wants=hermes-dashboard.service

[Service]
Type=simple
User=admin
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/admin/.hermes/hermes-agent/venv/bin/python3 /home/admin/hermes-bridge/agent_state_monitor.py --watch
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
